# docker/supervisor/worker-queue.conf
[program:account-queue-worker]
directory=/var/www/html
command=wait-for-it rabbitmq:5672 -- php yii queue/listen
autostart=true
autorestart=true
priority=10
