# docker/supervisor/worker-queue.conf
[program:queue-worker]
directory=/var/www/html
command=wait-for-it rabbitmq:5672 -- php yii queue/listen -v
autostart=true
autorestart=true
priority=10
